# scripts/pyro.service
[Unit]
Description=Pyro Robot Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pyro
WorkingDirectory=/home/pyro/pyro
Environment=REPO_DIR=/home/pyro/pyro
Environment=VENV_PATH=/home/pyro/venv
Environment=LOOP_HZ=150
Environment=KD=0.05
Environment=OUTPUT_EMA_ALPHA=0.8
Environment=PITCH_ALPHA=0.65
Environment=CAL_A=-4.76,2.95
Environment=CAL_B=-4.42,-10.70
Environment=CAL_C=7.37,-3.29
Environment=ROLL_OFFSET=-4.5
Environment=MAX_TILT_DEG=45
Environment=DEBUG=1
ExecStart=/home/pyro/pyro/scripts/start_pyro.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
